# scenarios/supervisord/wan-fed-gossip/scenario.ini
[program:consul-server-dc1]
command = consul agent 
  -server 
  -node server-%(process_num)01d
  -data-dir data/dc1/server-%(process_num)01d
  -bind 127.0.0.%(process_num)01d
  -client 127.0.0.%(process_num)01d
  -config-file config/dc1/server.hcl
  -config-dir config/dc1/server-%(process_num)01d
numprocs = 3
numprocs_start = 1
process_name = %(process_num)01d

[program:consul-client-dc1]
command = consul agent 
  -node client-%(process_num)01d
  -data-dir data/dc1/client-%(process_num)01d
  -bind 127.0.0.%(process_num)01d
  -client 127.0.0.%(process_num)01d
  -config-file config/dc1/client.hcl
  -config-dir config/dc1/client-%(process_num)01d
numprocs = 3
numprocs_start = 4
process_name = %(process_num)01d

[program:consul-server-dc2]
command = consul agent 
  -server 
  -node server-%(process_num)01d
  -data-dir data/dc2/server-%(process_num)01d
  -bind 127.0.0.%(process_num)01d
  -client 127.0.0.%(process_num)01d
  -config-file config/dc2/server.hcl
  -config-dir config/dc2/server-%(process_num)01d
numprocs = 3
numprocs_start = 7
process_name = %(process_num)01d

[program:consul-client-dc2]
command = consul agent 
  -node client-%(process_num)01d
  -data-dir data/dc2/client-%(process_num)01d
  -bind 127.0.0.%(process_num)01d
  -client 127.0.0.%(process_num)01d
  -config-file config/dc2/client.hcl
  -config-dir config/dc2/client-%(process_num)01d
numprocs = 3
numprocs_start = 10
process_name = %(process_num)01d
