# scenarios/supervisord/mesh-cluster-locality-aware/scenario.ini
[program:consul-server]
command = consul agent 
  -server 
  -node server-%(process_num)01d
  -data-dir data/server-%(process_num)01d
  -bind 127.0.0.%(process_num)01d
  -client 127.0.0.%(process_num)01d
  -config-file config/server.hcl
  -config-dir config/server-%(process_num)01d
numprocs = 3
numprocs_start = 1
process_name = %(process_num)01d

[program:consul-client]
command = consul agent 
  -node client-%(process_num)01d
  -data-dir data/client-%(process_num)01d
  -bind 127.0.0.%(process_num)01d
  -client 127.0.0.%(process_num)01d
  -config-file config/client.hcl
  -config-dir config/client-%(process_num)01d
numprocs = 3
numprocs_start = 4
process_name = %(process_num)01d

[program:frontend]
environment = BIND_ADDR=127.0.0.%(process_num)01d:6060,BACKEND_URL=http://127.0.0.%(process_num)01d:7000
command = frontend
numprocs = 1
numprocs_start = 4
process_name = %(process_num)01d

[program:backend]
environment = BIND_ADDR=127.0.0.%(process_num)01d:7070,BACKEND_HOSTNAME=client-%(process_num)01d
command = backend
numprocs = 2
numprocs_start = 5
process_name = %(process_num)01d

[program:envoy-frontend]
command = consul connect envoy
  -sidecar-for frontend 
  -admin-bind 127.0.0.%(process_num)01d:19000
  -http-addr http://127.0.0.%(process_num)01d:8500
  -grpc-addr http://127.0.0.%(process_num)01d:8502
  -token root 
  - --
  - -l debug
numprocs = 1
numprocs_start = 4
startretries=300
process_name = %(process_num)01d

[program:envoy-backend]
command = consul connect envoy
  -sidecar-for backend 
  -admin-bind 127.0.0.%(process_num)01d:19000
  -http-addr http://127.0.0.%(process_num)01d:8500
  -grpc-addr http://127.0.0.%(process_num)01d:8502
  -token root
  - --
  - -l debug
numprocs = 2
numprocs_start = 5
startretries=300
process_name = %(process_num)01d
