# systemd/envoy-dashboard.service
[Unit]
Description=Consul Service Mesh Envoy Proxy for Service Dashboard
After=network.target consul.service
Requires=consul.service
ConditionPathExists=/etc/consul.d/dashboard.hcl

[Service]
Type=simple
ExecStart=/usr/bin/consul connect envoy -sidecar-for=dashboard -- -l debug
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target

